# scripts/deploy/split-tunneling/dnsmasq.service.d/override.conf
# =============================================================================
# dnsmasq drop-in: зависимость от awg1 + автоперезапуск
#
# Кладётся в /etc/systemd/system/dnsmasq.service.d/override.conf
#
# Поскольку dnsmasq слушает на 10.9.0.1 (адрес интерфейса awg1), он не сможет
# запуститься, если awg1 ещё не поднят. Зависимость After/Wants гарантирует
# правильный порядок старта при загрузке системы.
#
# Restart=always — переживаем случайный crash или временное падение awg1.
# =============================================================================

[Unit]
After=awg-quick@awg1.service network-online.target
Wants=awg-quick@awg1.service network-online.target

[Service]
Restart=always
RestartSec=2
